# files/nftables-firewall/nftables-firewall.service
[Unit]
Description=nftables Firewall Control System
Documentation=man:nft(8)
DefaultDependencies=no
After=network-pre.target systemd-modules-load.service local-fs.target
Before=network.target shutdown.target
Conflicts=shutdown.target iptables.service ip6tables.service firewalld.service ufw.service
ConditionPathExists=/etc/nftables-firewall/

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nfwctl start
ExecReload=/usr/local/bin/nfwctl reload
ExecStop=/bin/true
RemainAfterExit=yes
TimeoutSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/nftables-firewall-backups /etc/nftables-firewall
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
